# .config/systemd/user/unlock-actions.service
[Unit]
Description=user actions on screen unlock

[Service]
ExecStart=%h/apps/unlock.sh

[Install]
WantedBy=unlock.target
